# recipes-connectivity/wifi/files/wifi_logger.service
[Unit]
Description=WiFi logger

[Service]
Type=simple
ExecStart=/usr/local/bin/wifi_logger.py
