# recipes-core/udev/udev/94-modems.rules
ACTION=="add",KERNEL=="ttyACM0",SUBSYSTEM=="tty",RUN+="/usr/local/bin/log_modem_detection.py &"
