# recipes-connectivity/cellular/files/cellular_data_supervisor.service
[Unit]
Description=Cellular Data Supervisor

[Service]
Type=simple
ExecStart=/usr/local/bin/cellular_data_supervisor.py

